# backend/investment-ledger.service
[Unit]
Description=Investment Ledger Backend
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/investment-pwa/backend
EnvironmentFile=/opt/investment-pwa/backend/.env
ExecStart=/opt/investment-pwa/backend/venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
